# systemd/hubado-restart.service
[Unit]
Description=Restart hubado (but only if it's running)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart hubado.service
